# systemd/user/wallpapers.service
[Unit]
Description=Random wallpapers

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wallpapers change
